# deploy/systemd/sisumail-update.timer
[Unit]
Description=Sisumail Update Timer (daily)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target

