# deploy/systemd/sisumail-relay.service
[Unit]
Description=Sisumail Relay API+SMTP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sisumail-relay.env
WorkingDirectory=/opt/sisumail
ExecStart=/usr/local/bin/sisumail-relay \
  -listen 127.0.0.1:8080 \
  -smtp-listen :25 \
  -smtp-zone sisumail.fi \
  -spool-ttl 15m \
  -purge-interval 30s \
  -db /var/lib/sisumail-relay.db \
  -max-json-bytes 1048576 \
  -max-ciphertext-bytes 262144 \
  -smtp-max-data-bytes 262144 \
  -smtp-max-connections 200 \
  -smtp-per-ip-per-minute 120 \
  -smtp-max-rcpt-per-session 10 \
  -admin-allow-cidrs ${SISUMAIL_ADMIN_ALLOW_CIDRS}
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
