# deploy/systemd/sisumail-relay.service
[Unit]
Description=Sisumail Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sisumail.env
ExecStart=/bin/sh -c '/usr/local/bin/sisumail-relay \
  -ssh-listen :2222 \
  -tier1-listen :2525 \
  -obs-listen "${SISUMAIL_OBS_LISTEN:-127.0.0.1:9090}" \
  -tier1-fast-fail-ms "${SISUMAIL_TIER1_FAST_FAIL_MS:-200}" \
  -tier1-open-timeout-ms "${SISUMAIL_TIER1_OPEN_TIMEOUT_MS:-3000}" \
  -tier1-idle-timeout-ms "${SISUMAIL_TIER1_IDLE_TIMEOUT_MS:-120000}" \
  -tier1-max-conn-duration-ms "${SISUMAIL_TIER1_MAX_CONN_DURATION_MS:-600000}" \
  -tier1-max-bytes-per-conn "${SISUMAIL_TIER1_MAX_BYTES_PER_CONN:-10485760}" \
  -tier1-max-conns-per-user "${SISUMAIL_TIER1_MAX_CONNS_PER_USER:-10}" \
  -tier1-max-conns-per-source "${SISUMAIL_TIER1_MAX_CONNS_PER_SOURCE:-20}" \
  -acme-dns01-per-user-per-min "${SISUMAIL_ACME_DNS01_PER_USER_PER_MIN:-30}" \
  -db /var/lib/sisumail/relay.db \
  -hostkey /var/lib/sisumail/hostkey_ed25519'
Restart=always
RestartSec=2

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sisumail /var/spool/sisumail

[Install]
WantedBy=multi-user.target
